# virtual/miller-base/files/daily-av-scan.cron
#!/bin/bash

freshclam --quiet
clamscan / -ri `mount | grep -v ^/dev | cut -d ' ' -f 3 | sed -e 's/^/--exclude-dir=/'` 
