# virtual/miller-base/files/sync-distclean-exec.cron
#!/bin/bash

eclean-dist -d -q -t 1w
eclean-pkg -d -q -t 1w
