# virtual/miller-base/files/sync-distclean-exec.cron
#!/bin/bash

eclean-dist -d -q
eclean-pkg -d -q
