# virtual/miller-base/files/update.portage.cron
#!/bin/bash

source "/etc/make.globals"
source "/etc/make.conf"

layman -SQ0
if df -l ${PORTDIR} 2> /dev/null > /dev/null ; then
    emerge --sync --quiet > /dev/null
fi

