# virtual/miller-base/files/update.local.portage.cron
#!/bin/bash

source "/sbin/functions.sh"
source "/etc/make.globals"
source "/etc/make.conf"

cd "${PORTDIR_OVERLAY}" && cvs -q up -dPCA
