# virtual/miller-base/files/update.local.portage.cron
#!/bin/bash

source "/etc/make.globals"
source "/etc/make.conf"

for d in ${PORTDIR_OVERLAY}; do
	if [ -d ${d}/CVS ] ; then
		cd "${d}" && cvs -q up -dPCA
	fi
done

