# virtual/miller-base/files/weekly-av-scan.cron
#!/bin/bash

eval $(emerge --info -v | grep '^[A-Z_]\+="[^"]*"$')

nice clamscan -ri --cross-fs=no \
	--exclude-dir=${PORTDIR} \
	--exclude-dir=${DISTDIR} \
	$(mount | grep ^/dev/ | sort -u -k1,1 | cut -d ' ' -f 3 | sort)
