# virtual/miller-base/files/weekly-av-scan.cron
#!/bin/bash

eval $(emerge --info -v | grep '^DISTDIR="[^"]*"$')

nice clamscan -ri --cross-fs=no \
	--exclude-dir=${DISTDIR} \
	$(emerge --info | sed -ne 's/\s\+location:\s\+\(.*\)/--exclude-dir="\1"/p') \
	$(mount | grep ^/dev/ | sort -u -k1,1 | cut -d ' ' -f 3 | sort)
