# virtual/miller-base/files/weekly-av-scan.cron
#!/bin/bash

source /etc/make.globals
source /etc/make.conf

nice clamscan / -ri ${PORTDIR} ${DISTDIR} `mount | grep -v ^/dev | grep -v ^rootfs | cut -d ' ' -f 3 | sed -e 's/^/--exclude-dir=/'` 
