# virtual/miller-base/files/weekly-distclean-exec-eclean.cron
#!/bin/bash

eclean-dist -f -d -q -t 2w
